# ansible-ec2-docker-swarm-tls/roles/certbot/templates/certbot-renewal.timer.j2
[Unit]
Description=Certbot Renewal Timer

[Timer]
OnCalendar={{ certbot_renewal_frequency }}
RandomizedDelaySec={{ certbot_renewal_delay }}
Persistent=true

[Install]
WantedBy=timers.target